# notifiers/wl-clipboard-notifier/wl-clipboard-notifier.service
[Unit]
Description=Wayland Clipboard Notifier

[Service]
ExecStart=%h/dotfiles/notifiers/wl-clipboard-notifier/wl-clipboard-notifier.sh

[Install]
WantedBy=default.target
