# notifiers/spotify-notifier/spotify-notifier.service
[Unit]
Description=Spotify Notifier

[Service]
ExecStart=%h/dotfiles/notifiers/spotify-notifier/spotify-notifier.sh

[Install]
WantedBy=default.target
